# build/sunshine.service
[Unit]
Description=Sunshine is a Gamestream host for Moonlight.

[Service]
ExecStart=sunshine

[Install]
WantedBy=graphical-session.target
